# install/units/metnos-stack-watchdog.timer.tmpl
[Unit]
Description=Periodic bounded readiness check for the Metnos stack
PartOf=metnos.target

[Timer]
OnActiveSec=3min
OnUnitActiveSec=2min
AccuracySec=15s
Unit=metnos-stack-watchdog.service

[Install]
WantedBy=metnos.target
